# rootfs/etc/udev/rules.d/91-usb_modeswitch.rules
SUBSYSTEM=="usb", RUN+="/etc/hotplug.d/usb/20-modeswitch"
