# rootfs/etc/hotplug.d/keypad/action
#!/bin/sh
RAKEY="/lib/wdk/rakey"

$RAKEY key${KEY}
