# rootfs/etc/init.d/gpio
#!/bin/sh /etc/rc.common
# Copyright (C) 2013 Montage Inc.
# OpenWrt.BOOT
START=06
start() {
	/bin/gpio > /dev/null
}
