# rootfs/etc/init.d/ledgpio
#!/bin/sh /etc/rc.common
# OpenWrt.BOOT
START=04
start() {
	/bin/led_gpio > /dev/null
}
