# rootfs/etc/init.d/umount
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
# OpenWrt.BOOT

STOP=99
stop() {
	sync
	umount -a -r
}
