# rootfs/usr/lib/opkg/info/rootfs.list
/bin/login.sh
/bin/ipcalc.sh
/bin/gpio
/bin/cdbak
/bin/pwr
/bin/led_gpio
/bin/madc
/etc/init.d/access_usb_storage.sh
/etc/init.d/alsa
/etc/init.d/alsa_record
/etc/init.d/ap_sta_bss_test
/etc/init.d/beacon_test
/etc/init.d/beacon_test_encrypt
/etc/init.d/beacon_test_n
/etc/init.d/beacon_test_n_encrpt
/etc/init.d/boot
/etc/init.d/br_test
/etc/init.d/br_test_n
/etc/init.d/br_test_n_encrpt
/etc/init.d/brgprobe
/etc/init.d/cdb
/etc/init.d/crond
/etc/init.d/done
/etc/init.d/ethopen
/etc/init.d/ethprobe
/etc/init.d/flash_test
/etc/init.d/fw
/etc/init.d/gpio
/etc/init.d/hostapd_cli
/etc/init.d/hwsim
/etc/init.d/hwsim_encrypt
/etc/init.d/iperf
/etc/init.d/lan
/etc/init.d/ledgpio
/etc/init.d/nat
/etc/init.d/ocfg
/etc/init.d/qa_20n_ccmp_wpa
/etc/init.d/qa_20n_ccmp_wpa2
/etc/init.d/qa_20n_tkip_wpa
/etc/init.d/qa_20n_tkip_wpa2
/etc/init.d/qa_20n_wep104
/etc/init.d/qa_20n_wep40
/etc/init.d/qa_40n_ccmp_wpa
/etc/init.d/qa_40n_ccmp_wpa2
/etc/init.d/qa_40n_tkip_wpa
/etc/init.d/qa_40n_tkip_wpa2
/etc/init.d/qa_40n_wep104
/etc/init.d/qa_40n_wep40
/etc/init.d/qa_b_ccmp_wpa
/etc/init.d/qa_b_ccmp_wpa2
/etc/init.d/qa_b_tkip_wpa
/etc/init.d/qa_b_tkip_wpa2
/etc/init.d/qa_b_wep104
/etc/init.d/qa_b_wep40
/etc/init.d/qa_g_ccmp_wpa
/etc/init.d/qa_g_ccmp_wpa2
/etc/init.d/qa_g_tkip_wpa
/etc/init.d/qa_g_tkip_wpa2
/etc/init.d/qa_g_wep104
/etc/init.d/qa_g_wep40
/etc/init.d/qos
/etc/init.d/rcS
/etc/init.d/sdio_usb
/etc/init.d/smp_test_eth_usb_sdio_flash_uart_test
/etc/init.d/station_test
/etc/init.d/station_test_encrpt
/etc/init.d/sysctl
/etc/init.d/telnetd
/etc/init.d/umount
/etc/init.d/usb
/etc/init.d/wan
/etc/init.d/wan_serv
/etc/init.d/watchdog
/etc/init.d/wlan
/etc/init.d/workmode
/etc/ifplugd/ifplugd.action
/etc/ppp/peers/Generic_disconn.scr
/etc/ppp/peers/Generic_conn.scr
/etc/udev/rules.d/91-usb_modeswitch.rules
/etc/hotplug.d/tty/30-3g_modem
/etc/hotplug.d/net/30-usbnic
/etc/hotplug.d/volume/action
/etc/hotplug.d/sound/soundcard
/etc/hotplug.d/usb/cdc-wdm
/etc/hotplug.d/usb/10-f_storage
/etc/hotplug.d/block/99-sysupgrade
/etc/hotplug.d/block/90-storage
/etc/hotplug.d/button/action
/etc/hotplug.d/keypad/action
/etc/config/modem.lst
/etc/config/cdb.dep
/etc/config/oom_disable.lst
/etc/config/apnlist.lst
/etc/protocols
/etc/group
/etc/dhcpcd.sh
/etc/passwd
/etc/preinit
/etc/hotplug2-common.rules
/etc/shells
/etc/rc.common
/etc/openwrt_release
/etc/profile
/etc/madconf
/etc/functions.sh
/etc/led_gpioconf
/etc/rc.local
/etc/hotplug2-init.rules
/etc/hosts
/etc/openwrt_version
/etc/gpioconf
/etc/sysctl.conf
/etc/diag.sh
/etc/asound.conf
/etc/inittab
/etc/banner
/etc/mtab
/etc/resolv.conf
/etc/fstab
/etc/TZ
/lib/upgrade/common.sh
/lib/functions/boot.sh
/lib/network/config.sh
/lib/preinit/10_essential_fs
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/40_run_failsafe_hook
/lib/preinit/90_restore_config
/lib/preinit/99_10_failsafe_login
/lib/preinit/05_mount_skip
/lib/preinit/30_device_fs_daemons
/lib/preinit/90_mount_no_jffs2
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/60_init_hotplug
/lib/preinit/10_indicate_failsafe
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/lib/preinit/50_indicate_regular_preinit
/lib/preinit/70_initramfs_test
/lib/preinit/99_10_run_init
/lib/preinit/10_indicate_preinit
/lib/preinit/40_init_shm
/lib/preinit/10_check_for_mtd
/lib/preinit/20_check_jffs2_ready
/lib/preinit/40_mount_jffs2
/lib/preinit/30_failsafe_wait
/lib/preinit/50_choose_console
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/lib/firstboot/20_reset_clear_jffs
/lib/firstboot/05_firstboot_skip
/lib/firstboot/10_no_fo_clear_overlay
/lib/firstboot/20_no_fo_mount_jffs
/lib/firstboot/10_reset_has_mini_fo
/lib/firstboot/30_is_rootfs_mounted
/lib/firstboot/50_pivot
/lib/firstboot/99_10_no_fo_cleanup
/lib/firstboot/20_has_mini_fo
/lib/firstboot/40_no_fo_copy_ramoverlay
/lib/firstboot/30_no_fo_pivot
/lib/firstboot/10_determine_parts
/lib/firstboot/99_10_with_fo_cleanup
/lib/firstboot/30_reset_copy_rom
/lib/firstboot/40_copy_ramoverlay
/rom/note
/sbin/sysupgrade
/sbin/wt
/sbin/wifi
/sbin/rf
/sbin/wd
/sbin/led.sh
/sbin/ifdown
/sbin/sw
/sbin/hotplug-call
/sbin/mount_root
/sbin/dw
/sbin/ew
/sbin/bb
/sbin/mdio
/sbin/hnat
/sbin/ifup
/sbin/clk
/sbin/firstboot
/usr/share/alsa/init/00main
/usr/share/alsa/init/cta0
/usr/share/alsa/init/cta1
/usr/share/udhcpc/default.script
/usr/lib/common.awk
/init
/www/udisk
/var
